# hybrid-shutdown-post.target
[Unit]
Description=Isolates System into Hibernate Mode
Requires=default.target
After=default.target
JobTimeoutSec=30
